# server/DEBIAN/prerm
#!/bin/sh -e
systemctl disable discover.service
systemctl stop discover.service
systemctl daemon-reload
